# ansible/roles/_common/templates/etc/profile.d/env-ansible.sh.j2
#
# {{ ansible_managed }}
#

function env-ansible {
  source {{ _common_ansible_venv_path }}/bin/activate
}
