# ansible/roles/_common/templates/etc/profile.d/env-ansible.sh.j2
#
# {{ ansible_managed }}
#

alias env-ansible='source /usr/local/bin/env-ansible.sh'
